# .github/deploy/lawroad-backend.service
[Unit]
Description=LawNRoad Backend Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/law-n-road-release
EnvironmentFile=/etc/lawroad/backend.env
ExecStart=/usr/bin/java -jar /opt/law-n-road-release/build/libs/law-n-road-0.0.1-SNAPSHOT.jar

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
