# .github/deploy/lawroad-backend.service
[Unit]
Description=LawNRoad Backend Service
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/lawroad/backend.env
WorkingDirectory=/opt/law-n-road-release
ExecStart=/usr/bin/java -jar /opt/law-n-road-release/libs/law-n-road-0.0.1-SNAPSHOT.jar
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=lawroad-backend

[Install]
WantedBy=multi-user.target
